# systemd/lightserver.socket
[Socket]
ListenStream=/run/lightserver.socket

[Install]
WantedBy=sockets.target
