# systemd/lightserver.service
[Unit]
Description=Example server of light sensor values

[Service]
ExecStart=/usr/sbin/lightserver
